# admin/package-files/systemd/flocker-container-agent.service
[Unit]
Description=Flocker Container Agent
After=docker.service
Wants=docker.service

[Service]
ExecStart=/usr/sbin/flocker-container-agent --journald
Restart=always

[Install]
WantedBy=multi-user.target
