# admin/package-files/systemd/flocker-dataset-agent.service
[Unit]
Description=Flocker Dataset Agent

[Service]
ExecStart=/usr/sbin/flocker-dataset-agent
Restart=always

# We can't have a private mount namespace, since we need to create mounts visible to Docker.

[Install]
WantedBy=multi-user.target
