# admin/package-files/systemd/flocker-control-api.socket
[Unit]
Description=Flocker Control REST API

[Socket]
ListenStream=4523
Service=flocker-control.service

[Install]
WantedBy=sockets.target
