# admin/package-files/systemd/flocker-dataset-agent.service
[Unit]
Description=Flocker Dataset Agent
AssertFileNotEmpty=/etc/flocker/dataset-agent.yml

[Service]
ExecStart=/usr/sbin/flocker-dataset-agent

# We can't have a private mount namespace, since we need to see ZFS created mounts.

[Install]
WantedBy=multi-user.target
