# admin/package-files/systemd/flocker-container-agent.service
[Unit]
Description=Flocker Container Agent
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/usr/sbin/flocker-container-agent
Restart=always

[Install]
WantedBy=multi-user.target
