# admin/package-files/systemd/flocker-control-agent.socket
[Unit]
Description=Flocker Control Agent Socket

[Socket]
ListenStream=4524
Service=flocker-control.service

[Install]
WantedBy=sockets.target
