# admin/package-files/flocker-control-agent.socket
[Unit]
Description=Flocker Control Agent Scoket

[Socket]
ListenStream=4524
NonBlocking=true

[Install]
WantedBy=sockets.target
